# nightColor.timer
[Unit]
Description=Run every day at 21:00 and 06:00.

[Timer]
OnCalendar=*-*-* *:01:00
Persistent=true

[Install]
WantedBy=timers.target
